# system-files/etc/systemd/system/systemd-networkd-wait-online.service.d/override.conf
# override timeout to be 1 second
[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-networkd-wait-online --timeout=1